# py-fastapi.service
[Unit]
Description=demo fastapi application
After=network.target

[Service]
User=osboxes
Group=osboxes
WorkingDirectory=/home/osboxes/app/src/
Environment="PATH=/home/osboxes/app/venv/bin"
EnvironmentFile=/home/osboxes/.env
ExecStart=/home/osboxes/app/venv/bin/gunicorn -w 2 -k uvicorn.workers.UvicornWorker app.main:app --bind 0.0.0.0:8000

[Install]
WantedBy=multi-user.target